# root-fs/etc/service.d/.s6-svscan/finish
#! /bin/bash
# INFRA HEADER -- DO NOT EDIT. metadata: {"checksum": "13857cba92955da681d86735f36ee26d"}

>&2 echo "[.s6-svscan/finish] [hostname:$(hostname)] Finished."