# docker/supervisor-common.conf
[program:splunk]
command = /opt/splunkforwarder/bin/splunk start splunkd --nodaemon
autorestart = true

[program:rsyslog]
command = rsyslogd -n
autorestart = true
priority = 1
