# docker/supervisor-web.conf
[program:uwsgi]
command = uwsgi --ini /root/uwsgi.ini
killasgroup = true

[program:nginx]
command = /usr/sbin/nginx
